# rootfs/etc/s6-overlay/s6-rc.d/unbound/run
#!/bin/sh

s6-echo "Starting unbound"

exec 2>&1
exec unbound -d
