# rootfs/etc/cron.d/cloudsync
*/15 * * * * root /usr/local/bin/cloudsync-cron
